# crontab
# 每分钟执行
* * * * * sh /weapp/web/weapp/services/service_minutely.sh &> /weapp/web/log/services/minutely.log

